# fluid-measurement.service
[Unit]
Description=Fluid Measurement Web Service
After=network.target pigpiod.service
Wants=pigpiod.service

[Service]
Type=simple
ExecStart=/opt/fluid_measurement/fluid_measurement_server
WorkingDirectory=/opt/fluid_measurement
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10
# Run as root for GPIO access with pigpio
User=root

[Install]
WantedBy=multi-user.target